# experimance/infra/systemd/audio@.service
[Unit]
Description=Audio Service (%i project)
After=experimance@%i.target
Wants=experimance@%i.target
PartOf=experimance@%i.target

[Service]
Type=simple
User=experimance
Group=experimance
WorkingDirectory=/opt/experimance
Environment=PROJECT_ENV=%i
Environment=EXPERIMANCE_ENV=production
Environment=PATH=/home/experimance/.pyenv/bin:/home/experimance/.local/bin:/usr/local/bin:/usr/bin:/bin
# Audio and D-Bus environment
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
# Wait for user session before starting audio (JACK needs user session)
ExecStartPre=/opt/experimance/infra/scripts/wait_for_desktop_session.sh experimance
ExecStart=/home/experimance/.local/bin/uv run -m experimance_audio
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=experimance-audio-%i

# Security
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/experimance /var/cache/experimance /home/experimance/.local /home/experimance/.pyenv /tmp /run/user/1000

# Audio access
SupplementaryGroups=audio

[Install]
WantedBy=experimance@.target
