# experimance/infra/systemd/core@.service
[Unit]
Description=Core Service (%i project)
After=network.target user-1000.slice
Wants=network.target
# Wait for experimance user session to be active (needed for RealSense camera access)
Requisite=user-1000.slice

[Service]
Type=simple
User=experimance
Group=experimance
WorkingDirectory=/opt/experimance
Environment=PROJECT_ENV=%i
Environment=EXPERIMANCE_ENV=production
Environment=PATH=/home/experimance/.pyenv/bin:/home/experimance/.local/bin:/usr/local/bin:/usr/bin:/bin
# Wait for user session before starting core service (camera access)
ExecStartPre=/opt/experimance/infra/scripts/wait_for_desktop_session.sh experimance
ExecStart=/home/experimance/.local/bin/uv run -m experimance_core
ExecReload=/bin/kill -HUP $MAINPID
Restart=no

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=experimance-core-%i

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/experimance /var/cache/experimance /home/experimance/.local /home/experimance/.pyenv

# Audio and video access for camera
SupplementaryGroups=video

[Install]
WantedBy=experimance@.target
