# experimance/infra/systemd/core@.service
# SYSTEMD TEMPLATE SERVICE
# This is a template file (note the @ in filename: core@.service)
# 
# Template Usage:
#   - This file defines HOW to run the core service
#   - %i gets replaced with the instance identifier (project name)
#   - To start: systemctl start core@experimance.service
#   - systemd creates instance "core@experimance.service" from this template
#   - Multiple projects can use same template: core@fire.service, etc.
#
# Instance Examples:
#   - core@experimance.service  (for experimance project)
#   - core@fire.service         (for fire project)
#   - core@demo.service         (for demo project)

[Unit]
Description=Core Service (%i project)
After=experimance@%i.target graphical.target
Wants=experimance@%i.target graphical.target
PartOf=experimance@%i.target

[Service]
Type=simple
User=experimance
Group=experimance
WorkingDirectory=/opt/experimance
Environment=PROJECT_ENV=%i
Environment=EXPERIMANCE_ENV=production
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PATH=/home/experimance/.pyenv/bin:/home/experimance/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/experimance/.local/bin/uv run -m experimance_core
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=experimance-core-%i

# Security
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/experimance /var/cache/experimance /var/log/experimance /home/experimance/.local /home/experimance/.pyenv /tmp /run/user/1000

# Audio and video access
SupplementaryGroups=audio video

[Install]
WantedBy=experimance@.target
