# experimance/infra/systemd/core@.service
# SYSTEMD TEMPLATE SERVICE
# This is a template file (note the @ in filename: core@.service)
# 
# Template Usage:
#   - This file defines HOW to run the core service
#   - %i gets replaced with the instance identifier (project name)
#   - To start: systemctl start core@experimance.service
#   - systemd creates instance "core@experimance.service" from this template
#   - Multiple projects can use same template: core@sohkepayin.service, etc.
#
# Instance Examples:
#   - core@experimance.service  (for experimance project)
#   - core@sohkepayin.service   (for sohkepayin project)
#   - core@demo.service         (for demo project)

[Unit]
Description=Core Service (%i project)
After=experimance@%i.target
Wants=experimance@%i.target
PartOf=experimance@%i.target

[Service]
Type=simple
User=experimance
Group=experimance
WorkingDirectory=/opt/experimance
Environment=PROJECT_ENV=%i
Environment=EXPERIMANCE_ENV=production
Environment=PATH=/home/experimance/.pyenv/bin:/home/experimance/.local/bin:/usr/local/bin:/usr/bin:/bin
# Wait for user session before starting core service (camera access)
ExecStartPre=/opt/experimance/infra/scripts/wait_for_desktop_session.sh experimance
ExecStart=/home/experimance/.local/bin/uv run -m experimance_core
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=experimance-core-%i

# Security
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/experimance /var/cache/experimance /home/experimance/.local /home/experimance/.pyenv /tmp

# Audio and video access for camera
SupplementaryGroups=video

[Install]
WantedBy=experimance@.target
