# experimance/infra/systemd/display@.service
[Unit]
Description=Display Service (%i project)
After=experimance@%i.target graphical-session.target
Wants=experimance@%i.target graphical-session.target
PartOf=experimance@%i.target

[Service]
Type=simple
User=experimance
Group=experimance
WorkingDirectory=/opt/experimance
Environment=PROJECT_ENV=%i
Environment=EXPERIMANCE_ENV=production
Environment=PATH=/home/experimance/.pyenv/bin:/home/experimance/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
ExecStart=/opt/experimance/infra/scripts/setup_display_env.sh /home/experimance/.local/bin/uv run -m experimance_display
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=experimance-display-%i

# Security
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/experimance /var/cache/experimance /home/experimance/.local /home/experimance/.pyenv /tmp

[Install]
WantedBy=experimance@.target
