# experimance/infra/systemd/experimance@.target
[Unit]
Description=Experimance Target (%i project)
After=network.target user@1000.service graphical-session.target
Wants=network.target user@1000.service graphical-session.target
Requires=network.target user@1000.service graphical-session.target
# Target ensures fundamental dependencies are met before any services start
# Services bind to this target and will start after it's ready

[Install]
WantedBy=multi-user.target
