# experimance/infra/systemd/experimance@.target
[Unit]
Description=Experimance Target (%i project)
After=network.target
Wants=network.target

[Install]
WantedBy=multi-user.target
