# experimance/infra/systemd/experimance-display@.user.service
[Unit]
Description=Experimance Display Service (%i project)
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
WorkingDirectory=/opt/experimance
Environment=EXPERIMANCE_PROJECT=%i
Environment=EXPERIMANCE_ENV=production
Environment=PATH=/home/experimance/.pyenv/bin:/home/experimance/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/experimance/.local/bin/uv run -m experimance_display
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=experimance-display-%i

# Resource limits
MemoryMax=2G
CPUQuota=100%

# Security - User services don't need as strict security
NoNewPrivileges=true

[Install]
WantedBy=default.target
