# experimance/infra/systemd/experimance-display@.service
[Unit]
Description=Experimance Display Service (%i project)
After=network.target experimance-core@%i.service
Wants=network.target
Requires=experimance-core@%i.service

[Service]
Type=simple
User=experimance
Group=experimance
WorkingDirectory=/home/experimance/experimance
Environment=EXPERIMANCE_PROJECT=%i
Environment=EXPERIMANCE_ENV=production
Environment=DISPLAY=:0
Environment=PATH=/home/experimance/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/experimance/.local/bin/uv run -m experimance_display
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=experimance-display-%i

# Resource limits
MemoryMax=2G
CPUQuota=100%

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/experimance/experimance/logs /home/experimance/experimance/images /var/cache/experimance

[Install]
WantedBy=experimance.target
