# experimance/infra/systemd/experimance-agent@.service
[Unit]
Description=Experimance Agent Service (%i project)
After=network.target
Wants=network.target

[Service]
Type=simple
User=experimance
Group=experimance
WorkingDirectory=/home/experimance/experimance
Environment=EXPERIMANCE_PROJECT=%i
Environment=EXPERIMANCE_ENV=production
Environment=PATH=/home/experimance/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/experimance/.local/bin/uv run -m experimance_agent
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=experimance-agent-%i

# Resource limits
MemoryMax=2G
CPUQuota=75%

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/experimance/experimance/logs /home/experimance/experimance/cache /var/cache/experimance /home/experimance/experimance/transcripts /home/experimance/experimance/media

# Audio and video access for microphone and camera
SupplementaryGroups=audio video

[Install]
WantedBy=experimance@%i.target
