# experimance/infra/systemd/experimance@.target
[Unit]
Description=Experimance Target (%i project)
After=network.target user-1000.slice
Wants=network.target
# Target waits for user session since all services now require it

[Install]
WantedBy=multi-user.target
