# experimance/infra/systemd/health@.service
[Unit]
Description=Experimance Health Service (%i)
After=network.target
Wants=network.target

[Service]
Type=simple
User=experimance
Group=experimance
WorkingDirectory=/opt/experimance
Environment="PROJECT_ENV=%i"
Environment="EXPERIMANCE_ENV=production"
Environment=PATH=/home/experimance/.pyenv/bin:/home/experimance/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/experimance/.local/bin/uv run -m experimance_health
Restart=no
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/experimance /var/cache/experimance /home/experimance/.local /home/experimance/.pyenv

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=experimance-health-%i

[Install]
WantedBy=multi-user.target
