# experimance/infra/systemd/matter-scheduler@.service
[Unit]
Description=Matter Device Scheduler for %i project
After=network.target
Wants=network.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/Documents/Projects/Art/experimance/installation/software/experimance
Environment=PATH=/home/%i/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/%i/Documents/Projects/Art/experimance/installation/software/experimance/libs/common/src

# Use uv to run the scheduler
ExecStart=/home/%i/.local/bin/uv run python scripts/matter_scheduler.py projects/fire/matter_schedule.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=matter-scheduler-%i

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/%i/Documents/Projects/Art/experimance/installation/software/experimance/logs

[Install]
WantedBy=default.target