# experimance/infra/systemd/reset-on-input.service
[Unit]
Description=Experimance Reset on Input - Keyboard/Controller Listener
After=multi-user.target
# Independent service - does not start or depend on experimance services

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/experimance
ExecStart=/usr/bin/python3 /opt/experimance/infra/scripts/reset_on_input.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Logging
SyslogIdentifier=experimance-reset-on-input

# Allow access to input devices
SupplementaryGroups=input

# Security
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/experimance /var/log/experimance

[Install]
WantedBy=multi-user.target
