# experimance/infra/launchd/com.experimance.fire.health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.experimance.health.fire</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/fireproject/Documents/experimance/experimance/infra/launchd/launch_health.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/fireproject/Documents/experimance/experimance</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/fireproject/Documents/experimance/experimance/logs/fire_health_launchd.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/fireproject/Documents/experimance/experimance/logs/fire_health_launchd_error.log</string>
    
    <!-- LaunchAgent version - no UserName needed, runs as logged-in user -->
    
    <!-- Restart on failure after 10 seconds -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <!-- Start after a delay to ensure system is ready -->
    <key>StartInterval</key>
    <integer>5</integer>
</dict>
</plist>
